# seo-api.service
[Unit]
Description=SEO API Server
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/mapsparser-Replit-front
Environment=PYTHONPATH=/root/mapsparser-Replit-front/src
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/mapsparser-Replit-front/venv/bin/python /root/mapsparser-Replit-front/src/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
