# seo-worker.service
[Unit]
Description=SEO Worker Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/mapsparser-Replit-front
Environment=PATH=/root/mapsparser-Replit-front/venv/bin
Environment=SUPABASE_URL=https://bvhpvzcvcuswiozhyqlk.supabase.co
Environment=SUPABASE_KEY=eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImJ2aHB2emN2Y3Vzd2lvemh5cWxrIiwicm9sZSI6InNlcnZpY2Vfcm9sZSIsImlhdCI6MTc1MjQ5OTg1OSwiZXhwIjoyMDY4MDc1ODU5fQ.NyT80y9DusZVhKD0Np36OQEYJPBjJK97UrTmJh60dlk
ExecStart=/root/mapsparser-Replit-front/venv/bin/python /root/mapsparser-Replit-front/src/worker.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 