# telegram-bot.service
[Unit]
Description=BeautyBot Telegram Bot
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/seo-app
Environment=PYTHONPATH=/opt/seo-app/src
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/seo-app/telegram-bot-venv/bin/python /opt/seo-app/src/telegram_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
