# ops/systemd/localos-storage-maintenance.timer
[Unit]
Description=Run LocalOS storage retention daily

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
